# meshnarc.service
[Unit]
Description=meshnarc — Meshtastic MQTT to BigQuery capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/meshnarc
ExecStart=/usr/bin/python3 /home/pi/meshnarc/meshnarc_sub.py
Restart=always
RestartSec=10

# Environment — edit these or use an env file
Environment=GOOGLE_APPLICATION_CREDENTIALS=/home/pi/meshnarc/sa-key.json
Environment=MESHNARC_PROJECT=your-gcp-project
Environment=MESHNARC_BROKER=your-broker.example.com
Environment=MESHNARC_MQTT_USER=meshnarc
Environment=MESHNARC_MQTT_PASS=changeme
Environment=MESHNARC_CAPTURE_LAT=0.0000
Environment=MESHNARC_CAPTURE_LON=0.0000

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/pi/meshnarc
PrivateTmp=true

[Install]
WantedBy=multi-user.target
